# scripts/launchd/com.chump.main-preflight-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Static reference copy. The authoritative plist is GENERATED by
     scripts/setup/install-main-preflight-watchdog-daemon-launchd.sh, which expands
     $ROOT, $HOME, and log paths at install time. This file exists for review /
     diff purposes only; do not load it directly. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chump.main-preflight-watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/Chump/scripts/coord/main-preflight-watchdog-daemon.sh</string>
        <string>tick</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/Chump</string>
    <key>StartInterval</key>
    <integer>600</integer>
    <!-- RunAtLoad=true: exercises the daemon immediately on install so
         "did the plist land correctly?" is answered at install time, not
         600 seconds later (INFRA-351 lesson). -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/home/.chump/logs/main-preflight-watchdog.out</string>
    <key>StandardErrorPath</key>
    <string>/path/to/home/.chump/logs/main-preflight-watchdog.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/path/to/home</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <!-- META-248: explicit absolute path to ambient.jsonl so the daemon
             does not compute it relative to a stale /tmp worktree under
             launchd's execution context. -->
        <key>CHUMP_AMBIENT_PATH</key>
        <string>/path/to/Chump/.chump-locks/ambient.jsonl</string>
    </dict>
</dict>
</plist>
