# scripts/launchd/com.chump.pr-shepherd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Static reference copy — META-192.
     The authoritative plist is GENERATED by install-pr-shepherd-daemon.sh,
     which expands $ROOT, $HOME, and log paths at install time.
     This file exists for review / diff purposes only; do not load it directly. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chump.pr-shepherd</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/Chump/scripts/coord/pr-shepherd-daemon.sh</string>
        <string>tick</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/Chump</string>
    <key>StartInterval</key>
    <integer>60</integer>
    <!-- RunAtLoad=true: exercises the daemon immediately on install so
         "did the plist land correctly?" is answered at install time, not
         60 seconds later (INFRA-351 lesson). -->
    <key>RunAtLoad</key>
    <true/>
    <!-- KeepAlive=false: each tick is a single-shot script exit.
         launchd re-launches it every StartInterval seconds. -->
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/path/to/home/.chump/logs/pr-shepherd.out</string>
    <key>StandardErrorPath</key>
    <string>/path/to/home/.chump/logs/pr-shepherd.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/path/to/home</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
