# scripts/launchd/com.chump.daemon-exit-loop-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Static reference copy. The authoritative plist is GENERATED by
     scripts/setup/install-daemon-exit-loop-watcher-launchd.sh, which expands
     $ROOT, $HOME, and log paths at install time. This file exists for review /
     diff purposes only; do not load it directly.
     INFRA-2417: STABLE PATH LESSON — daemon must use /Users/<user>/Projects/Chump,
     NOT /private/tmp/... The integrator-daemon crash-loop root cause (2026-06-02)
     was a /private/tmp/chump-install dir that vanished. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chump.daemon-exit-loop-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeffadkins/Projects/Chump/scripts/coord/daemon-exit-loop-watcher-daemon.sh</string>
        <string>tick</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jeffadkins/Projects/Chump</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <!-- RunAtLoad=true: exercises the daemon immediately on install so
         "did the plist land correctly?" is answered at install time, not
         900 seconds later (INFRA-351 lesson). -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jeffadkins/.chump/logs/daemon-exit-loop-watcher.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffadkins/.chump/logs/daemon-exit-loop-watcher.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jeffadkins</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <!-- META-248: explicit absolute path to ambient.jsonl so the daemon
             does not compute it relative to a stale /tmp worktree under
             launchd's execution context. -->
        <key>CHUMP_AMBIENT_PATH</key>
        <string>/Users/jeffadkins/Projects/Chump/.chump-locks/ambient.jsonl</string>
        <key>REPO_ROOT</key>
        <string>/Users/jeffadkins/Projects/Chump</string>
    </dict>
</dict>
</plist>
