# scripts/launchd/com.chump.ghost-pr-closer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Static reference copy — META-225.
     The authoritative plist is GENERATED by install-ghost-pr-closer.sh,
     which expands $ROOT, $HOME, and log paths at install time.
     This file exists for review / diff purposes only; do not load it directly. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chump.ghost-pr-closer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/Chump/scripts/coord/ghost-pr-closer.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/Chump</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/home/.chump/logs/ghost-pr-closer.out</string>
    <key>StandardErrorPath</key>
    <string>/path/to/home/.chump/logs/ghost-pr-closer.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/path/to/home</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
