# scripts/launchd/com.chump.trunk-sentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Static reference copy. The authoritative plist is GENERATED by
     scripts/setup/install-trunk-sentinel.sh, which expands $ROOT, $HOME,
     and log paths at install time. This file exists for review / diff
     purposes only; do not load it directly. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chump.trunk-sentinel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/Chump/scripts/coord/trunk-sentinel-daemon.sh</string>
        <string>tick</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/Chump</string>
    <key>StartInterval</key>
    <integer>60</integer>
    <!-- RunAtLoad=true: exercises the daemon immediately on install so
         "did the plist land correctly?" is answered at install time, not
         60 seconds later (INFRA-351 lesson). -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/home/.chump/logs/trunk-sentinel.out</string>
    <key>StandardErrorPath</key>
    <string>/path/to/home/.chump/logs/trunk-sentinel.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/path/to/home</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <!-- INFRA-2340: auth placeholders. The installer (install-trunk-sentinel.sh)
             substitutes these from ~/.chump/oauth-token.json + env ANTHROPIC_API_KEY,
             or deletes the lines entirely when both are absent. -->
        <key>CLAUDE_CODE_OAUTH_TOKEN</key>
        <string>CLAUDE_CODE_OAUTH_TOKEN_PLACEHOLDER</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>ANTHROPIC_API_KEY_PLACEHOLDER</string>
    </dict>
</dict>
</plist>
