# scripts/plists/dev.chump.queue-driver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.chump.queue-driver</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/jeffadkins/Projects/Chump &amp;&amp; bash scripts/coord/queue-driver.sh</string>
  </array>
  <!-- INFRA-1304: every 5 min. queue-driver does:
       - INFRA-1081: BEHIND PR scan from cache, gh pr update-branch refresh
       - INFRA-1137: DIRTY PR scan, rebase + merge-driver auto-resolve
       Closes the last manual gap in the PR drain. -->
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/chump-queue-driver.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/chump-queue-driver.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/jeffadkins/Projects/Chump</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- gh CLI needs HOME to find the auth token. -->
    <key>HOME</key>
    <string>/Users/jeffadkins</string>
    <key>PATH</key>
    <string>/Users/jeffadkins/.local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
