# scripts/plists/dev.chump.lint-storm-detector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.chump.lint-storm-detector</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/jeffadkins/Projects/Chump/scripts/ops/lint-storm-detector.sh</string>
  </array>
  <!-- Run once per hour. Detects when 3+ open PRs fail the same clippy lint
       and files an INFRA gap to relax it at workspace.lints.clippy (INFRA-672). -->
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/chump-lint-storm-detector.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/chump-lint-storm-detector.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/jeffadkins/Projects/Chump</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- gh CLI needs HOME to find the auth token. -->
    <key>HOME</key>
    <string>/Users/jeffadkins</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:$HOME/.cargo/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
