# scripts/plists/dev.chump.sccache-reaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.chump.sccache-reaper</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/jeffadkins/Projects/Chump/scripts/coord/sccache-reaper.sh --execute</string>
  </array>
  <!-- Run every 6 hours. Checks ~/Library/Caches/Mozilla.sccache size and
       prunes oldest files when over SCCACHE_CACHE_CAP_GB (default 10GB).
       Emits kind=sccache_reaped to .chump-locks/ambient.jsonl. INFRA-2303. -->
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/chump-sccache-reaper.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/chump-sccache-reaper.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/jeffadkins/Projects/Chump</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/jeffadkins</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/jeffadkins/.cargo/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
