# scripts/plists/doc-keeper.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.chump.doc-keeper</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/path/to/Chump/scripts/dev/doc-keeper.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>
  <key>RunAtLoad</key>
  <true/>
  <!-- 21600 = 6 h. Use 86400 for daily if preferred. -->
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>StandardOutPath</key>
  <string>/path/to/Chump/logs/doc-keeper-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Chump/logs/doc-keeper-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/you/.local/bin</string>
  </dict>
</dict>
</plist>
