# scripts/plists/heartbeat-health-check.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.chump.heartbeat-health-check</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/path/to/Chump/scripts/ci/check-heartbeat-health.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>1200</integer>
  <!-- 1200 = 20 min. Checks heartbeat logs and recommends interval adjustments for peak performance. -->
  <key>StandardOutPath</key>
  <string>/path/to/Chump/logs/heartbeat-health.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Chump/logs/heartbeat-health.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin</string>
    <key>CHUMP_HOME</key>
    <string>/path/to/Chump</string>
  </dict>
</dict>
</plist>
