# scripts/plists/hourly-update-to-discord.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.chump.hourly-update-to-discord</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/path/to/Chump/scripts/dev/hourly-update-to-discord.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>
  <key>RunAtLoad</key>
  <false/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <!-- 3600 = 1 hour. Sends hourly summary DM via notify tool. Requires CHUMP_READY_DM_USER_ID and DISCORD_TOKEN in .env. -->
  <key>StandardOutPath</key>
  <string>/path/to/Chump/logs/hourly-update-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Chump/logs/hourly-update-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/you/.local/bin</string>
  </dict>
</dict>
</plist>
