# scripts/plists/oven-tender.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.chump.oven-tender</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/path/to/Chump/scripts/dev/oven-tender.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <!-- 3600 = 1 hour. Pre-warms model (warm-the-ovens) if not already warm. For Ollama set WARM_PORT=11434 in .env or here. -->
  <key>StandardOutPath</key>
  <string>/path/to/Chump/logs/oven-tender-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Chump/logs/oven-tender-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/you/.local/bin</string>
    <!-- Max mode default: WARM_PORT=8000 (vLLM). If using Ollama: <key>WARM_PORT</key><string>11434</string> -->
  </dict>
</dict>
</plist>
