# scripts/setup/com.chump.mesh-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- INFRA-2545 (FLEET-034 Phase 2): the work-routing mesh CONSUMER daemon.
     Pairs with com.chump.coord-assign (the publisher). Picks one capability-
     matched, mesh-routed gap per tick. ROUTE-ONLY by default (proves routing
     is live, zero execution risk); set CHUMP_MESH_WORKER_EXECUTE=1 below to
     hand routed gaps to `chump --execute-gap`. Fail-open: NATS down -> exit 0,
     pull fleet handles work. Install:
       cp scripts/setup/com.chump.mesh-worker.plist ~/Library/LaunchAgents/
       launchctl load -w ~/Library/LaunchAgents/com.chump.mesh-worker.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.chump.mesh-worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/jeffadkins/Projects/Chump &amp;&amp; bash scripts/coord/mesh-worker-loop.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CHUMP_NATS_URL</key><string>nats://127.0.0.1:4222</string>
    <key>WORKER_SKILLS</key><string>rust,sqlite,coord,git,bash,macos,docs</string>
    <key>WORKER_MACHINE</key><string>macbook</string>
    <key>WORKER_BACKEND</key><string>claude</string>
    <!-- Route-only by default. Flip to 1 to enable autonomous execution: -->
    <key>CHUMP_MESH_WORKER_EXECUTE</key><string>0</string>
  </dict>
  <key>WorkingDirectory</key><string>/Users/jeffadkins/Projects/Chump</string>
  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>300</integer>
  <key>StandardErrorPath</key><string>/Users/jeffadkins/.chump/logs/mesh-worker.err</string>
  <key>StandardOutPath</key><string>/Users/jeffadkins/.chump/logs/mesh-worker.out</string>
</dict>
</plist>
