# scripts/com.chump.battle-qa.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chump.battle-qa</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/jeffadkins/Projects/Chump &amp;&amp; ./scripts/run-battle-qa-full.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/jeffadkins/Projects/Chump</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/jeffadkins/Projects/Chump/logs/battle-qa-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jeffadkins/Projects/Chump/logs/battle-qa-launchd.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
