# scripts/heartbeat-shepherd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.chump.heartbeat-shepherd</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/path/to/Chump/scripts/heartbeat-shepherd.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>900</integer>
  <!-- 900 = 15 min. Checks heartbeat-learn.log; optionally retries one quick round if last failed. -->
  <key>StandardOutPath</key>
  <string>/path/to/Chump/logs/heartbeat-shepherd-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Chump/logs/heartbeat-shepherd-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin</string>
    <!-- <key>HEARTBEAT_SHEPHERD_RETRY</key><string>1</string> -->
  </dict>
</dict>
</plist>
