# scripts/restart-vllm-if-down.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.chump.restart-vllm-if-down</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/path/to/Chump/scripts/restart-vllm-if-down.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>180</integer>
  <!-- 180 = 3 min. Keeps vLLM-MLX on 8000 up when Python/MLX crashes. Only starts if 8000 is down. -->
  <key>StandardOutPath</key>
  <string>/path/to/Chump/logs/restart-vllm-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Chump/logs/restart-vllm-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/you/.local/bin</string>
  </dict>
</dict>
</plist>
