# scripts/shed-load.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.chump.shed-load</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/path/to/Chump/scripts/enter-chump-mode.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>
  <key>RunAtLoad</key>
  <false/>
  <key>StartInterval</key>
  <integer>7200</integer>
  <!-- 7200 = 2 hours. Stops Ollama + embed and kills blocklisted apps (scripts/chump-mode.conf) so vLLM-MLX + Chump have max GPU/RAM. Edit chump-mode.conf to choose what to quit; run list-heavy-processes.sh to see candidates. -->
  <key>StandardOutPath</key>
  <string>/path/to/Chump/logs/shed-load-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/Chump/logs/shed-load-launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/Users/you/.local/bin</string>
  </dict>
</dict>
</plist>
