# launchd/com.chump.opus-curator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  com.chump.opus-curator — INFRA-842 LaunchAgent.

  Runs scripts/coord/opus-curator.sh every 10 minutes (600s).
  Proactive health audit: P0 inflation, vague gaps, pillar balance,
  stale PRs, waste rate. Primary fleet decision-maker.

  INSTALLATION (use scripts/setup/install-curator-launchd.sh):
    1. scripts/setup/install-curator-launchd.sh
    2. Verify: launchctl list | grep chump.opus-curator
    3. Manual run: launchctl start com.chump.opus-curator
       Logs: tail /tmp/chump-opus-curator.out.log

  UNINSTALL:
    launchctl unload ~/Library/LaunchAgents/com.chump.opus-curator.plist
    rm ~/Library/LaunchAgents/com.chump.opus-curator.plist

  Tunables (env vars read by opus-curator.sh):
    CHUMP_CURATOR_DRY_RUN=1   audit only, no mutations
    CHUMP_CURATOR_VERBOSE=1   extra logging
-->
<dict>
  <key>Label</key>
  <string>com.chump.opus-curator</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>scripts/coord/opus-curator.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/chump-opus-curator.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/chump-opus-curator.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/path/to/HOME</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/path/to/.cargo/bin:/usr/bin:/bin</string>
    <!-- INFRA-2310: consensus merge gate SHADOW mode (log-only, never blocks).
         Observe kind=consensus_gate_would_block in ambient for ~7 days,
         then operator sets =enforce (per INFRA-2274 / META-195). -->
    <key>CHUMP_CONSENSUS_MERGE_GATE</key>
    <string>1</string>
  </dict>

  <key>ThrottleInterval</key>
  <integer>120</integer>
</dict>
</plist>
