# launchd/com.chump.bounced-pr-detector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  com.chump.bounced-pr-detector — INFRA-781 LaunchAgent.

  Runs scripts/coord/bounced-pr-detector.sh every 30 min. Detects PRs
  that closed without merging where the work has not re-landed via
  another PR; emits ambient ALERTs and auto-files recovery gaps so
  bounced work doesn't disappear.

  Origin story: 2026-05-08, PR #1349 (EVAL-026 default-flip) was
  closed by an agent when its branch got recycled for unrelated work.
  The default-flip never made it to main. Without this watcher, that
  kind of bounce is invisible.

  INSTALLATION (manual; use scripts/setup/install-bounced-pr-detector-launchd.sh
  for the templated path-replace + launchctl load):
    1. Edit this file: replace /path/to/Chump with your repo path.
    2. cp launchd/com.chump.bounced-pr-detector.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.chump.bounced-pr-detector.plist
    4. Verify: launchctl list | grep chump.bounced-pr-detector
    5. Manual run: launchctl start com.chump.bounced-pr-detector
       Logs: tail /tmp/chump-bounced-pr-detector.out.log

  UNINSTALL:
    launchctl unload ~/Library/LaunchAgents/com.chump.bounced-pr-detector.plist
    rm ~/Library/LaunchAgents/com.chump.bounced-pr-detector.plist

  Tunables (set as EnvironmentVariables below or override at runtime):
    CHUMP_BOUNCED_PR_DETECTOR=0           — disable entirely
    CHUMP_BOUNCED_PR_LOOKBACK_HOURS=24    — how far back to scan
    CHUMP_BOUNCED_PR_AUTO_FILE_GAP=1      — auto-file recovery gaps (default on)
-->
<dict>
  <key>Label</key>
  <string>com.chump.bounced-pr-detector</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/path/to/Chump/scripts/coord/bounced-pr-detector.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/Chump</string>

  <!-- Every 30 min: same cadence as ghost-gap-reaper. -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/chump-bounced-pr-detector.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/chump-bounced-pr-detector.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/opt/homebrew/sbin</string>
    <!-- gh + git + python3 must be on PATH. -->
  </dict>

  <key>ThrottleInterval</key>
  <integer>60</integer>
</dict>
</plist>
